# scripts/bootstrap/ubuntu24/build-agent.service
[Unit]
Description=Build Agent Peer Discovery
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/build-agent/agent-runner.sh
Restart=on-failure
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=300
TimeoutStopSec=30
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
